# crontab
0 0,12,18 * * * cd code/influxdb-backup && ./cronjob
0 6       * * * cd code/influxdb-backup && ./cronjob --full
